# luci-app-rakitanmanager/root/etc/uci-defaults/99_rakitanmanager
#!/bin/bash
# Copyright 2024 RTA SERVER

echo "Setup php uhttpd"
uci set uhttpd.main.index_page='index.php'
uci set uhttpd.main.interpreter='.php=/usr/bin/php-cgi'
uci commit uhttpd

/etc/init.d/uhttpd restart
echo "Setup php uhttpd Done"

echo "Setup Package For Python3"
if which pip3 >/dev/null; then
    # Instal paket 'requests' jika belum terinstal
    if ! pip3 show requests >/dev/null; then
        echo "Installing package 'requests'"
        if ! pip3 install requests >/dev/null; then
            echo "Error installing package 'requests'"
            echo "Setup Gagal | Mohon Coba Kembali"
            exit 1  # Keluar dari skrip dengan status error
        fi
    else
        echo "Package 'requests' already installed"
    fi

    # Instal paket 'huawei-lte-api' jika belum terinstal
    if ! pip3 show huawei-lte-api >/dev/null; then
        echo "Installing package 'huawei-lte-api'"
        if ! pip3 install huawei-lte-api >/dev/null; then
            echo "Error installing package 'huawei-lte-api'"
            echo "Setup Gagal | Mohon Coba Kembali"
            exit 1  # Keluar dari skrip dengan status error
        fi
    else
        echo "Package 'huawei-lte-api' already installed"
    fi

    # Instal paket 'datetime' jika belum terinstal
    if ! pip3 show datetime >/dev/null; then
        echo "Installing package 'datetime'"
        if ! pip3 install datetime >/dev/null; then
            echo "Error installing package 'datetime'"
            echo "Setup Gagal | Mohon Coba Kembali"
            exit 1  # Keluar dari skrip dengan status error
        fi
    else
        echo "Package 'datetime' already installed"
    fi

    # Instal paket 'logging' jika belum terinstal
    if ! pip3 show logging >/dev/null; then
        echo "Installing package 'logging'"
        if ! pip3 install logging >/dev/null; then
            echo "Error installing package 'logging'"
            echo "Setup Gagal | Mohon Coba Kembali"
            exit 1  # Keluar dari skrip dengan status error
        fi
    else
        echo "Package 'logging' already installed"
    fi
else
    echo "Error: 'pip3' command not found"
    echo "Setup Gagal | Mohon Coba Kembali"
    exit 1  # Keluar dari skrip dengan status error
fi

sed -i 's/\r$//' /usr/share/rakitanmanager/core-manager.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/modem-hilink.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/modem-mf90.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/modem-hp.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/modem-rakitan.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/plugins/adb-deviceinfo.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/plugins/adb-sms.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/plugins/service-openclash.sh
sed -i 's/\r$//' /usr/share/rakitanmanager/plugins/systeminfo.sh

echo "Setup ModemManager"
mm1="/usr/lib/ModemManager/connection.d/10-report-down"
mm2="/usr/lib/ModemManager/connection.d/10-report-down-and-reconnect"
mm3="/usr/lib/ModemManager/connection.d/rakitanmanager"

if [ -f "$mm1" ]; then
    rm /usr/lib/ModemManager/connection.d/10-report-down
fi
if [ -f "$mm3" ]; then
    if [ -f "$mm2" ]; then
        rm /usr/lib/ModemManager/connection.d/10-report-down-and-reconnect
    fi
    mv "/usr/lib/ModemManager/connection.d/rakitanmanager" "/usr/lib/ModemManager/connection.d/10-report-down"
    chmod +x /usr/lib/ModemManager/connection.d/10-report-down
fi

#Restore
if [ -f "/tmp/rakitanmanager.bak" ]; then
	mv -f "/tmp/rakitanmanager.bak" "/etc/config/rakitanmanager" >/dev/null 2>&1
    mv -f "/tmp/data-modem.json.bak" "/usr/share/rakitanmanager/data-modem.json" >/dev/null 2>&1
    mv -f "/tmp/bot_message.txt.bak" "/www/rakitanmanager/bot_message.txt" >/dev/null 2>&1
	rm -rf "/tmp/rakitanmanager.bak" >/dev/null 2>&1
	rm -rf "/tmp/data-modem.json.bak" >/dev/null 2>&1
	rm -rf "/tmp/bot_message.txt.bak" >/dev/null 2>&1
fi
uci set rakitanmanager.cfg.enabled='0' && uci commit rakitanmanager
exit
